# private_dot_local/bin/scripts/ledger-web.service
[Unit]
Description=Ledger web service

[Service]
ExecStart=/root/ledger.sh
StartLimitIntervalSec=10sec
Restart=on-failure

[Install]
WantedBy=default.target
